# supervisor-conf/smonitor.conf
[program:smonitor]
user=yukai
command=/home/yukai/.venvs/smonitor/bin/python smonitor.py
directory=/home/yukai/projects/smonitor
autostart=true
autorestart=true
